# vm-conf/etc/systemd/system/netctl@bridge.vm.service
.include /usr/lib/systemd/system/netctl@.service

[Unit]
Description=VM Bridge
BindsTo=sys-subsystem-net-devices-br0.device
